# network/wifi-watchdog.service
[Unit]
Description=WiFi connectivity watchdog (escape to known-good SSID if no internet)
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/wifi-watchdog.sh
